# configs/services/hyprmo_desktop_widget.service
[Unit]
Description=Hyprmo: start desktop widget (ie. clock)

[Service]
Type=simple
Restart=always
ExecStart=hyprmo_hook_desktop_widget.sh
