# configs/services/hyprmo_hook_lisgd.service
[Unit]
Description=run user defined lisgd start hook

[Service]
Type=simple
Restart=always
ExecStart=hyprmo_hook_lisgdstart.sh
